# build/rfid-reader.service
[Unit]
Description=RFID Reader
After=multi-user.target

[Service]
Type=idle
ExecStart=/opt/escaperoom/bin/rfid-reader.py
User=pi
Environment=DISPLAY=:0
Restart=always

[Install]
WantedBy=multi-user.target
